# dot_config/autostart/gnome-keyring-ssh.desktop
[Desktop Entry]
Type=Application
Name=SSH Key Agent
Exec=/usr/bin/true
Hidden=true
